# deb/DEBIAN/postinst
#!/bin/sh -e

case "$1" in
    configure|reconfigure)
        systemctl daemon-reload
        systemctl enable hd-idle
    ;;
esac